# roles/firewall/files/firewall.target
[Unit]
Description=Firewall target

[Install]
WantedBy=multi-user.target
