# roles/python2/files/uwsgi@.socket
[Unit]
Description=uWSGI socket for %i

[Socket]
ListenStream=/run/uwsgi/%i.socket
User=http
Group=http

[Install]
WantedBy=sockets.target

