# roles/python2/files/uwsgi@.socket
[Unit]
Description=uWSGI socket for %i

[Socket]
ListenStream=/run/uwsgi/%i.socket
User=%i
Group=%i

[Install]
WantedBy=sockets.target

